# source/pbmssm/build/deb/prerm
#!/bin/bash
# bmssm deb prerm：卸载/升级前停服务。
# $1: remove/upgrade/deconfigure
set -e

systemctl stop bmssm 2>/dev/null || true

# 仅卸载/清除时禁用并删 unit；升级保留（postinst 会重启新版本）
if [ "$1" = "remove" ] || [ "$1" = "purge" ]; then
  systemctl disable bmssm 2>/dev/null || true
  [ "$1" = "purge" ] && rm -f /etc/systemd/system/bmssm.service
fi

exit 0
